# etc/init.d/vmtoolsd
#!/sbin/openrc-run

command="/usr/sbin/vmwtools"
pidfile="/run/$(RC_SVCNAME).pid"
command_args="--background=$(pidfile)"

